# bahram-cm/deploy/supervisor/bahram-family-queue.conf
; Bahram CM — Family ("خانواده داداش بهرام") queue workers
; Isolated from the default bahram-queue pool so heavy media transfer/processing
; jobs never starve regular order/notification jobs (or vice versa).
; Copy to /etc/supervisor/conf.d/bahram-family-queue.conf
; Then: sudo supervisorctl reread && sudo supervisorctl update

[program:bahram-family-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/bahram-cm/backend/artisan queue:work redis --queue=family-high,family-notifications,family-ai,family-analytics,family-media,family-low --sleep=3 --tries=3 --max-time=3600 --timeout=600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/bahram-family-queue.log
stopwaitsecs=600
