# bahram-cm/deploy/supervisor/bahram-queue.conf
; Bahram CM — Laravel queue workers
; Copy to /etc/supervisor/conf.d/bahram-queue.conf
; Then: sudo supervisorctl reread && sudo supervisorctl update

[program:bahram-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/bahram-cm/backend/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --timeout=120
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/bahram-queue.log
stopwaitsecs=3600
